# claw-runner.service
[Unit]
Description=claw-runner (KRunner integration for OpenClaw / Clawdbot)
After=graphical-session.target

[Service]
Type=simple
# Make PATH predictable for resolving the configured CLI (e.g. ~/.local/bin/clawdbot)
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/giles/claw-runner/bin/claw-runner
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
